# deploy/supervisor/rastreador-workerman.ini
; ============================================================
; Supervisor — Listener TCP Assíncrono (Workerman)
; CentOS/RHEL: copiar para /etc/supervisord.d/rastreador-workerman.ini
; ============================================================
[program:rastreador-workerman]
command=/usr/bin/php /var/www/rastreador/artisan workerman:listen start
directory=/var/www/rastreador
user=nginx
autostart=true
autorestart=true
startretries=5
startsecs=5
redirect_stderr=true
stdout_logfile=/var/log/supervisor/rastreador-workerman.log
stdout_logfile_maxbytes=50MB
stdout_logfile_backups=5
environment=APP_ENV="production",HOME="/var/www/rastreador",USER="nginx"
